# agents/elitedesk/elitedesk-watcher.service
[Unit]
Description=home-ops elitedesk watcher
After=network.target docker.service

[Service]
Type=simple
EnvironmentFile=%h/.config/elitedesk-watcher.env
ExecStart=/usr/bin/python3 %h/bin/elitedesk-watcher.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
