# ops/elitedesk/daily-digest-evening.service
[Unit]
Description=home-ops evening digest email + queue overnight LLM narrative
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/logs-stack
ExecStart=%h/logs-stack/ops/elitedesk/daily-digest.sh --mode evening
# Stream stdout/stderr to journald — elitedesk-watcher picks it up into host_logs.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
