# ops/elitedesk/pg-backup.service
[Unit]
Description=home-ops Postgres logical backup (pg_dump -Fc) to NAS
After=docker.service network-online.target remote-fs.target
Wants=remote-fs.target

[Service]
Type=oneshot
# Run from the repo clone; the script does its own cd-less work via $BACKUP_DIR.
WorkingDirectory=%h/logs-stack
ExecStart=%h/logs-stack/ops/elitedesk/pg-backup.sh
# Stream stdout/stderr to journald — gets picked up by elitedesk-watcher and lands in host_logs.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
